# php-websocket/supervisord.d/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan socket:serve
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
